# docker/backend/controller.service
[Unit]
Description=Foo

[Service]
ExecStart=/home/julia/backend/run.sh

[Install]
WantedBy=multi-user.target
